# sv/nexus/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars
env

cd /nexus/bin
exec 2>&1
exec ./nexus console
